# Python/old/celery/celery_demo/conf/supervisor_celery_flower.ini
[program:celery_demo-celery-flower]
command=python manage.py celery flower -l INFO
directory=/media/file/project/python/celery/celery_demo
environment=PATH="/home/zyp/.pyenv/versions/Python3.6/bin"
stdout_logfile=/media/file/project/python/celery/celery_demo/logs/celery.flower.log
stderr_logfile=/media/file/project/python/celery/celery_demo/logs/celery.flower.log

autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
priority=1000
