# Python/old/celery/celery_demo/conf/supervisor_celery_worker.ini
[program:celery_demo-celery-worker]
command=python manage.py celery worker -l INFO
directory=/media/file/project/python/celery/celery_demo
environment=PATH="/home/zyp/.pyenv/versions/Python3.6/bin"
stdout_logfile=/media/file/project/python/celery/celery_demo/logs/celery.worker.log
stderr_logfile=/media/file/project/python/celery/celery_demo/logs/celery.worker.log

autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
priority=998
